# src/ndt_slam/scripts/ndt-slam.service
[Unit]
Description=NDT-SLAM Long-Term Mapping
After=network.target

[Service]
Type=simple
User=ydkj
WorkingDirectory=/home/ydkj/NDT-slam-ws
ExecStartPre=/bin/bash -c "mkdir -p /home/ydkj/slam_data/maps/live/current"
ExecStart=/bin/bash -c "source /home/ydkj/NDT-slam-ws/devel/setup.bash && roslaunch ndt_slam warehouse_live_longterm_mapping.launch"
Restart=on-failure
RestartSec=15
StartLimitIntervalSec=300
StartLimitBurst=5
# 防止子进程残留
KillMode=mixed
TimeoutStopSec=10

[Install]
WantedBy=multi-user.target
